# authelia/rootfs/etc/s6-overlay/s6-rc.d/authelia/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant App: authelia
# Runs authelia
# ==============================================================================

bashio::log.info 'Starting authelia...'

# Check if authelia binary exists
if [ -f /usr/bin/authelia ]; then
    bashio::log.debug 'Authelia binary found!'
else
    bashio::log.error 'Authelia binary not found at /usr/bin/authelia'
    exit 1
fi

# Check if the configuration file exists
if [ -f /etc/authelia/config.template.yml ]; then
    bashio::log.debug 'Configuration file found!'
else
    bashio::log.error 'Configuration file not found at /etc/authelia/config.template.yml'
    exit 1
fi

# Run authelia
bashio::log.debug 'Running authelia...'
exec /usr/bin/authelia --config /config/config.yml